# app.conf
[program:app]
command=node server.js
autostart=true