# app.conf
[program:app]
command=node server.js > log.txt
autostart=true
autorestart = true
startsecs=3